# lego/usr/lib/systemd/system/letsencrypt.timer
[Unit]
Description=Renew Let's Encrypt certificate

[Timer]
OnCalendar=monday *-*-* 05:00:00
Persistent=true
Unit=letsencrypt.service

[Install]
WantedBy=timers.target
